# scripts/nimos-daemon.service
[Unit]
Description=NimOS Daemon
Documentation=https://github.com/andresgv-beep/NimOs-beta-8.1
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/nimos/nimos.env
ExecStart=/opt/nimos/daemon/nimos-daemon
Restart=always
RestartSec=3
TimeoutStartSec=120

# ─── Security hardening (A3) ────────────────────────────────────────────────
# Reduces the blast radius of any RCE: a compromise becomes partial, not total.
# The daemon still needs mount/mkfs.btrfs/wipefs/chgrp + user management, so we
# bound capabilities instead of dropping root entirely, and whitelist the paths
# it must write. VALIDATE ON HARDWARE: if storage/shares/backup break, the most
# likely cause is a missing ReadWritePaths entry or a capability removed below.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
PrivateTmp=true
ReadWritePaths=/opt/nimos /var/lib/nimos /var/log/nimos /etc/nimos /etc/exports /etc/samba /mnt /srv /home
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
# Capabilities required: mount (CAP_SYS_ADMIN), ownership changes (CAP_CHOWN,
# CAP_DAC_OVERRIDE, CAP_FOWNER), user/identity management (CAP_SETUID, CAP_SETGID),
# binding privileged ports if ever needed (CAP_NET_BIND_SERVICE).
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_SETUID CAP_SETGID CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_SYS_ADMIN CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_SETUID CAP_SETGID CAP_NET_BIND_SERVICE

# Logging
StandardOutput=append:/var/log/nimos/daemon.log
StandardError=append:/var/log/nimos/daemon-error.log

[Install]
WantedBy=multi-user.target
